# ops/yatco-automation/systemd/moana-yatco-refresh.timer
[Unit]
Description=Run the Moana YATCO refresh every 72 hours

[Timer]
OnBootSec=20min
OnUnitInactiveSec=72h
AccuracySec=5min
RandomizedDelaySec=15min
Unit=moana-yatco-refresh.service

[Install]
WantedBy=timers.target
